# systemd/sec-lsm-manager.service.in
[Unit]
Description=Security Manager service

[Service]
ExecStart=/usr/bin/@CMAKE_PROJECT_NAME@d --user @USER@ --group @GROUP@ --groups @CYNAGORA_GROUP@

Type=notify

KillMode=process
TimeoutStopSec=3
Restart=always

Sockets=@CMAKE_PROJECT_NAME@.socket

CapabilityBoundingSet=CAP_MAC_ADMIN CAP_DAC_OVERRIDE CAP_MAC_OVERRIDE CAP_SETUID CAP_SETGID

#NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
