# systemd/sec-lsm-manager.target
[Unit]
Description=sec-lsm-manager target
DefaultDependencies=true

