# systemd/sec-lsm-manager.socket.in
[Unit]
Description=sec-lsm-manager socket

[Socket]
FileDescriptorName=sec-lsm-manager
ListenStream=@SOCKET_DIR@/sec-lsm-manager.socket
SocketUser=@USER@
SocketGroup=@GROUP@
SocketMode=0660
SmackLabelIPIn=@
SmackLabelIPOut=@

Service=sec-lsm-manager.service

[Unit]
Wants=sec-lsm-manager.target
Before=sec-lsm-manager.target

[Install]
WantedBy=sockets.target
