# systemd/security-manager.service.in
[Unit]
Description=Security Manager service

[Service]
ExecStart=/usr/bin/security-managerd --user @USER@ --group @GROUP@ --groups @CYNAGORA_GROUP@

Type=notify

KillMode=process
TimeoutStopSec=3
Restart=always

Sockets=security-manager.socket

CapabilityBoundingSet=CAP_MAC_ADMIN CAP_DAC_OVERRIDE CAP_MAC_OVERRIDE CAP_SETUID CAP_SETGID

#NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
