# systemd/security-manager.target
[Unit]
Description=security-manager target
DefaultDependencies=true

